# cron/cert-update.service
[Unit]
Description=SSL certificate renewal

[Service]
Type=simple
ExecStart=/home/ubuntu/renew.sh
