# cron/outdating.service
[Unit]
Description=Offer outdating routine

[Service]
Type=simple
ExecStart=/home/ubuntu/outdate.sh
