# cron/scraping.service
[Unit]
Description=Scraping routine

[Service]
Type=simple
ExecStart=/home/ubuntu/scrape.sh
